# supervisor.conf
[program:biyunxd]    
command=/home/ubuntu/pywebsites/biyunxd/venv/bin/uwsgi --ini /home/ubuntu/pywebsites/biyunxd/uwsgi.ini
